# .crontab

47 * * * * anacron -t /home/sam/.anacrontab -S /home/sam/.acancron/spool/ -s
